# etc/profile.d/date_time_format.sh
alias date='date +"%I:%M %p, %d/%m/%Y"'
