# etc/profile.d/date_time_format.sh
Set time format to 12-hour with AM/PM and date format to dd/mm/yyyy
